# salt/roots/salt/go-backend/go-backend.conf
[program:go-backend]
command=bash /src/run-go-backend.sh
user=vagrant
stopsignal=INT
autostart=true
autorestart=true
stopasgroup=true
startretries=999999999
