# salt/roots/salt/python-backend/python-backend.conf
[program:python-backend]
command=bash /src/run-python-backend.sh
user=vagrant
stopsignal=INT
autostart=true
autorestart=true
stopasgroup=true
startretries=999999999
