# salt/roots/salt/pusud/pusud.conf
[program:pusud]
command=bash /src/run-pusud.sh
user=vagrant
stopsignal=INT
autostart=true
autorestart=true
stopasgroup=true
startretries=999999999
